# centos7-slave-jenkins/filesystem_root/etc/profile.d/99_env_setup.sh
export JAVA_HOME=/usr/lib/jvm/jre-1.8.0-openjdk
export JAVACMD=/usr/lib/jvm/jre-1.8.0-openjdk/bin/java
export ANT_HOME=/opt/ant
pathmunge $ANT_HOME/bin
pathmunge $JAVA_HOME/bin
